# back-end/scripts/PlayerScripts/crontab
* * * * * root /usr/local/bin/python /app/DailyCronJob.py >> /var/log/cron.log 2>&1
